# systemd/system/log-temp.timer
[Unit]
Description=Temperature Log Timer

[Timer]
#OnBootSec=1min
#OnUnitActiveSec=15min
OnCalendar=*:00/10

[Install]
WantedBy=timers.target
